# pi-display/trakio-display.service
# systemd service — auto-start trakio kiosk on boot
#
# Install:
#   sudo cp trakio-display.service /etc/systemd/system/
#   sudo systemctl enable trakio-display
#   sudo systemctl start  trakio-display
#
# To pass a fixed WS URL instead of using the on-screen form:
#   ExecStart=/home/pi/trakio-display/start.sh ws://192.168.43.1:9001

[Unit]
Description=trakio Pi Kiosk Display
After=graphical-session.target network.target
Wants=graphical-session.target

[Service]
User=pi
Group=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
WorkingDirectory=/home/pi/trakio-display
ExecStart=/home/pi/trakio-display/start.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=graphical-session.target
